# crondns.cron
*/5 * * * * www-data php /var/www/crondns/lib/dyn_dns_update.php >> /var/www/crondns/data/cron.log 2>&1
